# nginx/runit/crond/run
#!/usr/bin/env sh

. /etc/envvars

exec /usr/sbin/crond -f -L /dev/stdout
