# nginx/runit/nginx/run
#!/usr/bin/env sh

. /etc/envvars

exec /usr/sbin/nginx
